# crontab
# Run the app every day at 2:15pm
15 * * * * root source /app/env_variables; PATH=/usr/local/bin:$PATH node /app/dist/index.js >> /var/log/cron.log 2>&1
